# solr4/deploy/supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock

[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:solr]
command = /bin/bash -c "cd /opt/solr/example; java -Xmx256m -jar start.jar"
priority = 999
stdout_logfile = /var/log/%(program_name)s.stdout
stderr_logfile = /var/log/%(program_name)s.stderr